# supervisord.conf
[supervisord]
nodaemon=true
[program:sshd]
command=/usr/sbin/sshd -D

[program:activemq]
command=/usr/local/apache-activemq-5.14.5/bin/activemq start
